# crontab/dump_to_jd.sh
#!/usr/bin/env bash
scrapy dump_to_jd &
