# crontab/dump_to_jd.sh
#!/usr/bin/env bash
scrapy mc_dump_to_jd &
