# crontab/dump_to_jd.sh
#!/usr/bin/env bash
cd /home/apps/proxy_spider
scrapy mc_dump_to_jd &
