# crontab/dump_to_jd.sh
#!/usr/bin/env bash

python ../proxy_spider/proxy/dump_to_jd.py >> ../log/dump_to_jd.log
