# crontab/dump_to_valid.sh
#!/usr/bin/env bash
scrapy dump_to_valid &