# crontab/scrawl.sh
#!/usr/bin/env bash
scrapy mc_crawlall &
