# crontab/valid_jd.sh
#!/usr/bin/env bash
scrapy valid_jd log/valid_jd.log