# crontab/valid_proxy_items_jd.sh
#!/usr/bin/env bash

python ../proxy_spider/proxy/valid_proxy_items_jd.py >> ../log/valid_proxy_items_jd.log