# crontab/valid_jd.sh
#!/usr/bin/env bash
scrapy mc_valid_jd &