# crontab/valid_valid.sh
#!/usr/bin/env bash
scrapy valid_valid log/valid_valid.log