# crontab/valid_valid.sh
#!/usr/bin/env bash
cd /home/apps/proxy_spider
scrapy mc_valid_valid